# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
user=vscode
autorestart=true
priority=100

[program:xfce]
command=/usr/bin/startxfce4
environment=DISPLAY=":1"
user=vscode
autorestart=true
priority=200
stopwaitsecs=5

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -nopw -rfbport 5901 -bg -o /var/log/x11vnc.log
user=vscode
autorestart=true
priority=300

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5901
user=vscode
autorestart=true
priority=400
